# tmux/launchd/com.einargudjonsson.tmux-prune.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.einargudjonsson.tmux-prune</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>/Users/einargudjonsson/dotfiles/tmux/scripts/tmux-prune.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/einargudjonsson/.cache/tmux-prune.launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/einargudjonsson/.cache/tmux-prune.launchd.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
